# machinegun.conf
start on startup
chdir /home/u/domains/machinegun
start on runlevel [2345]
exec /home/u/domains/machinegun/warp

